# autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=xx-daemon conky --make-pidfile --oknodo --start
