# systemd/user/emacs.service.d/daemon-name.conf
[Service]
ExecStart=
ExecStart=/usr/bin/emacs --fg-daemon=systemd
ExecStop=
ExecStop=/usr/bin/emacsclient --eval "(kill-emacs)" --server-file=systemd
