# systemd/user/emacs.service.d/timeout.conf
[Service]
# Sometimes Emacs daemon upgrades ELPA packages during start-up.  This may take
# a long time.  Disable TimeoutStartSec.
TimeoutStartSec=infinity
