# src/dot_config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=o daemon conky --background --make-pidfile --oknodo --start
