# src/dot_config/autostart/keynav.desktop
[Desktop Entry]
Type=Application
Name=keynav
Exec=o daemon keynav --background --make-pidfile --oknodo --start
