# src/dot_config/autostart/keynav.desktop
[Desktop Entry]
Type=Application
Name=keynav
Exec=xx-daemon keynav --background --make-pidfile --oknodo --start
