# src/dot_config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=nm-applet
Exec=o daemon nm-applet --background --make-pidfile --oknodo --start
