# src/dot_config/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Name=sxhkd
Exec=o daemon sxhkd --background --make-pidfile --oknodo --start
