# src/dot_config/autostart/unclutter.desktop
[Desktop Entry]
Type=Application
Name=unclutter
Exec=o daemon unclutter --background --make-pidfile --oknodo --start -- -idle 1 -root
