# src/dot_config/autostart/warpd.desktop
[Desktop Entry]
Type=Application
Name=warpd
Exec=o daemon warpd --background --make-pidfile --oknodo --start
