# roles/config/templates/etc/systemd/system/docker.service.d/env.conf
[Service]
Environment="DOCKER_OPTS=--disable-legacy-registry=true {{ config_docker_opts }}"
EnvironmentFile=-/etc/default/docker
