# roles/config/templates/etc/systemd/system/docker.service.d/exec.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --disable-legacy-registry=true {{ config_docker_opts }} $DOCKER_OPTS
