# systemd/user/conky.service
[Unit]
Description=conky

[Service]
Type=forking
ExecStart=/bin/sh -l -c 'conky'
Restart=on-failure
