# systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Service]
Type=forking
ExecStart={{ HOME }}/bin/dropbox.py start
ExecStop={{ HOME }}/bin/dropbox.py stop
