# systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Service]
Type=forking
ExecStart=/bin/sh -l -c 'dropbox.py start'
ExecStop=/bin/sh -l -c 'dropbox.py stop'
Restart=on-failure
