# systemd/user/gpg-agent.service
[Unit]
Description=gpg-agent

[Service]
Type=forking
ExecStart=/bin/sh -l -c 'gpg-agent --daemon'
Restart=on-failure
