# systemd/user/keynav.service
[Unit]
Description=keynav

[Service]
Type=forking
ExecStart=/usr/bin/keynav
