# systemd/user/keynav.service
[Unit]
Description=keynav

[Service]
Type=simple
ExecStart=/bin/sh -l -c 'exec keynav'
Restart=on-failure
