# systemd/user/keynav.service
[Unit]
Description=keynav

[Service]
Type=forking
ExecStart=/bin/sh -l -c 'keynav'
Restart=on-failure
