# systemd/user/xautolock.service
[Unit]
Description=xautolock

[Service]
ExecStart=/usr/bin/xautolock -time 20 -locker "xlock.sh --yes"
Restart=on-failure
