# systemd/user/xwrits.service
[Unit]
Description=xwrits

[Service]
ExecStart=/usr/bin/xwrits +multiscreen +breakclock +cheat +clock +idle +mouse +top
