# roles/config/files/etc/systemd/system/syncthing@.service.d/env.conf
[Service]
Environment=STNORESTART=yes STNOUPGRADE=yes
EnvironmentFile=-/etc/default/syncthing
